# sait/sait3/usr/lib/systemd/system/proxyfsd.service
[Unit]
Description=Swift ProxyFS service

[Service]
Environment=GOTRACEBACK=crash
LimitCORE=infinity
ExecStart=/vagrant/bin/proxyfsd /vagrant/src/github.com/swiftstack/ProxyFS/sait/sait3/proxyfs.conf
ExecReload=/usr/bin/kill -HUP $MAINPID
Restart=no
KillMode=process
